# infra/deploy/digitalocean/avalanche.service
[Unit]
Description=Avalanche homeserver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=avalanche
Group=avalanche
EnvironmentFile=/etc/avalanche/avalanche.env
ExecStart=/usr/local/lib/avalanche/avalanche-server
WorkingDirectory=/var/lib/avalanche
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/avalanche

[Install]
WantedBy=multi-user.target
